# symlinked/config/kitty/update-kitty.crontab
00 12 * * * root   $HOME/.config/kitty/update-kitty.sh
